# data/init/pritunl.conf
description "Pritunl Daemon"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn

exec pritunl --pidfile /var/run/pritunl.pid
