# data/init/pritunl.conf
description "Pritunl Daemon"

start on runlevel [2345]
stop on runlevel [!2345]

exec pritunl start
