# arch/git/pritunl.install
pre_remove() {
    systemctl stop pritunl || true
    systemctl disable pritunl || true
}

post_remove() {
    rm -rf /var/lib/pritunl
}
